# app/gunicorn.service
[Unit]
Description=demo fastapi application
After=network.target

[Install]
WantedBy=multi-user.target